# ivote-server/debian/evote-hes.cron.d
# monitor
*/1 * * * *	hes /usr/share/evote/hes_monitor.py

