# ivote-server/debian/evote-hts.postrm
#!/bin/sh

set -e

NAME=evote-hts
USER=hts

if [ "$1" = remove ] ; then
	rm -f /etc/apache2/sites-enabled/hts-apache
	rm -f /etc/apache2/sites-available/hts-apache
	invoke-rc.d apache2 restart || true
fi

if [ "$1" = purge ] ; then
	deluser --quiet $USER > /dev/null || true
fi


#DEBHELPER#

